# seagull/root/etc/services.d/seagull/run
#!/usr/bin/with-contenv sh

exec /go/src/github.com/tobegit3hub/seagull/seagull
