# kibana/root/etc/services.d/kibana/run
#!/usr/bin/env bash

echo "Starting Kibana"
exec /opt/kibana/bin/kibana;
