# consul/root/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/dnsmasq -d -C /etc/dnsmasq.conf;
