# nginx/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

[[ -f ./configure ]] && ./configure

exec nginx;
