# ubuntu-14.04/root/etc/services.d/serf/run
#!/usr/bin/with-contenv sh

USER=serf
CONFIG=/etc/serf
BIN=/usr/sbin/serf

exec s6-setuidgid $USER $BIN agent -config-dir $CONFIG;
