# centos-7/root/etc/services.d/logstash-forwarder/run
#!/usr/bin/env bash

[[ -f ./configure ]] && source configure

CONFIG=/etc/logstash-forwarder.conf
BIN=/opt/logstash-forwarder/bin/logstash-forwarder

exec $BIN -config $CONFIG;
