# chef-server/root/etc/services.d/chef-server/run
#!/usr/bin/with-contenv sh

exec /etc/services.d/chef-server/init.rb
