# centos-7/root/etc/services.d/syslog/run
#!/usr/bin/env bash

[[ -f ./configure ]] && source configure

exec rsyslogd -f /etc/rsyslog.conf -n
